# lib/itamae/plugin/recipe/letsencrypt/cron/templates/etc/cron.d/itamae-letsencrypt.erb
# DO NOT EDIT
# BECAUSE THIS CRON CREATE BY itamae-plugin-recipe-letsencrypt
0 0 1 * * <%= @user %> <%= @command %> renew
